# deploy/multiserver/supervisor-worker.conf
[program:revibe-queue-worker]
command=/usr/bin/php /var/www/revibe/scripts/queue_worker.php --daemon --limit=50
directory=/var/www/revibe
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=30
numprocs=1
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/revibe-queue-worker.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
environment=APP_ENV="production"
